# ui.service
# /etc/systemd/system/ui.service
[Unit]
Description=UI Frontend Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/AI_Email_Assistant
ExecStart=/usr/bin/npm run dev
Restart=always
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

